# monitor-avahi.service
[Unit]
Description=Monitor Avahi deamon for invalid hostnames
BindsTo=avahi-daemon.service

[Service]
ExecStart=/usr/bin/monitor-avahi
Restart=on-failure

[Install]
WantedBy=multi-user.target
